# pa_xtouch_control.service
[Unit]
Description=PulseAudio X-Touch control daemon

[Service]
ExecStart=/usr/bin/pa_xtouch_control.sh

[Install]
WantedBy=default.target