# pamidi.service
[Unit]
Description=Control PulseAudio with a Midi device

[Service]
ExecStart=/usr/bin/pamidi

[Install]
WantedBy=default.target
